# files/server/clustercheck.socket
[Unit]
Description=MySQL Clustercheck Socket for HAProxy

[Socket]
ListenStream=9200
ReusePort=true
Accept=true

[Install]
After=mysql.service
